# deploy/package/facturatie-backend.service
[Unit]
Description=Facturatie ERP Backend API
After=network.target mongod.service
Wants=mongod.service

[Service]
Type=simple
User=facturatie
Group=facturatie
WorkingDirectory=/home/facturatie/htdocs/facturatie.sr/backend
Environment="PATH=/home/facturatie/htdocs/facturatie.sr/backend/venv/bin"
ExecStart=/home/facturatie/htdocs/facturatie.sr/backend/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8001 --workers 2
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/facturatie/htdocs/facturatie.sr/logs/backend.log
StandardError=append:/home/facturatie/htdocs/facturatie.sr/logs/backend.error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
